# crontab/crontab_notice.sh
#!/bin/bash
cd /data/web/monitor/crontab/
/usr/local/php/bin/php crontab_analysisData.php
