# crontab/crontabtest_gendata.sh
#!/bin/bash
cd /data/web/monitor/service/php/api/
/usr/local/php/bin/php test_report.php
